# deploy/supervisord.conf
[supervisord]
nodaemon = true

[program: uwsgi]
command=uwsgi --ini subzero/deploy/uwsgi.ini

[program:nginx]
command=nginx -g "daemon off;" -c /etc/nginx/nginx.conf